# deploy/garuda-chatbot.service
[Unit]
Description=Garuda Food Chatbot
After=network.target

[Service]
User=root
WorkingDirectory=/var/www/garuda-chatbot
EnvironmentFile=/var/www/garuda-chatbot/.env
ExecStart=/var/www/garuda-chatbot/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:5000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
